# structr-binaries/src/main/resources/deb/control/postinst
#!/bin/sh
#
# Post installation script for structr [[version]]
#

systemctl enable structr